# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:idrivecron]
command=perl /etc/idrivecron.pl